# supervisor/tomcat.conf
[program:tomcat]
command=/usr/bin/start-tomcat.sh